# mlss-monitor.service
[Unit]
Description=MLSS Monitor Flask App
After=network.target

[Service]
# Edit User and WorkingDirectory to match your installation
User=masadmin
WorkingDirectory=/home/masadmin/projects/git_versions/mars-air-quility
ExecStart=/home/masadmin/.local/bin/poetry run python mlss_monitor/app.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Non-secret config (intervals, paths, plug IP) — readable by anyone
EnvironmentFile=%h/projects/git_versions/mars-air-quility/.env

# Secrets (SECRET_KEY, GitHub OAuth credentials) — root-owned, mode 600
# Create with: sudo mkdir -p /etc/mlss && sudo touch /etc/mlss/secrets.env
#              sudo chmod 600 /etc/mlss/secrets.env && sudo chown root:root /etc/mlss/secrets.env
# systemd reads this as root before dropping to User=masadmin, so masadmin
# never has direct read access to the file — only the injected env vars.
EnvironmentFile=/etc/mlss/secrets.env

[Install]
WantedBy=multi-user.target
